# addon/rootfs/etc/services.d/gsm_sms/run
#!/usr/bin/with-contenv bashio

# Wait for initialization to complete
sleep 2

bashio::log.info "Starting GSM SMS Service..."

# Load configuration
CONFIG_PATH=/data/options.json
MODEM_DEVICE=$(bashio::config 'device')
BAUD_SPEED=$(bashio::config 'baud_speed')
UNICODE=$(bashio::config 'unicode')
SCAN_INTERVAL=$(bashio::config 'scan_interval')
SERVICE_NAME=$(bashio::config 'service_name')

# Set up the GSM SMS service
bashio::log.info "Setting up Legacy GSM SMS service..."

# Create configuration for the SMS service
mkdir -p /data/gsm_sms
cat > /data/gsm_sms/config.json << EOL
{
  "device": "${MODEM_DEVICE}",
  "baud_speed": "${BAUD_SPEED}",
  "unicode": ${UNICODE},
  "scan_interval": ${SCAN_INTERVAL},
  "service_name": "${SERVICE_NAME}"
}
EOL

# Check if modem exists before starting service
if [ ! -e "$MODEM_DEVICE" ]; then
  bashio::log.warning "GSM modem device not found: $MODEM_DEVICE"
  bashio::log.info "Available devices:"
  ls -la /dev/tty*
  
  # We'll continue anyway, as the device might appear later
  bashio::log.info "Will continue trying to connect..."
fi

bashio::log.info "Starting GSM SMS Service with device: $MODEM_DEVICE"

# Execute the Python service directly with proper environment
cd /app
export PYTHONUNBUFFERED=1
export S6_KEEP_ENV=1
exec python3 /app/gsm_sms_service.py
