# addon/rootfs/etc/services.d/gsm_sms/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Starting GSM SMS Service..."

# Execute script
exec /run.sh
