# addon-v2/rootfs/etc/services.d/gsm-sms/finish
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Finish/cleanup script for GSM SMS service
# ==============================================================================

bashio::log.info "GSM SMS service stopped"
