# z-addon-old-test/rootfs/etc/s6-overlay/s6-rc.d/init-gsm-config/up
#!/command/with-contenv bashio
# S6 oneshot service to create the GSM SMS configuration

# Notify S6 when we're done (this is important for oneshot services)
S6_NOTIFICATION_FD=3
if [ -f "/var/run/s6/container_environment/S6_NOTIFICATION_FD" ]; then
    S6_NOTIFICATION_FD=$(cat /var/run/s6/container_environment/S6_NOTIFICATION_FD)
fi

bashio::log.info "Running GSM configuration setup (S6_NOTIFICATION_FD=${S6_NOTIFICATION_FD})..."

# Create configuration file for GSM SMS service
CONFIG_PATH=/data/options.json
MODEM_DEVICE=$(bashio::config 'device')
BAUD_SPEED=$(bashio::config 'baud_speed')
UNICODE=$(bashio::config 'unicode')
SCAN_INTERVAL=$(bashio::config 'scan_interval')
SERVICE_NAME=$(bashio::config 'service_name')
LOG_LEVEL=$(bashio::config 'log_level')

# Debug logging
bashio::log.debug "Creating GSM SMS configuration in /data/gsm_sms/config.json"
bashio::log.debug "Device: ${MODEM_DEVICE}"
bashio::log.debug "Baud Speed: ${BAUD_SPEED}"
bashio::log.debug "Unicode: ${UNICODE}"
bashio::log.debug "Scan Interval: ${SCAN_INTERVAL}"
bashio::log.debug "Service Name: ${SERVICE_NAME}"
bashio::log.debug "Log Level: ${LOG_LEVEL}"

# Create configuration for the SMS service
mkdir -p /data/gsm_sms
cat > /data/gsm_sms/config.json << EOL
{
  "device": "${MODEM_DEVICE}",
  "baud_speed": "${BAUD_SPEED}",
  "unicode": ${UNICODE},
  "scan_interval": ${SCAN_INTERVAL},
  "service_name": "${SERVICE_NAME}",
  "log_level": "${LOG_LEVEL}"
}
EOL

# Verify file was created successfully
if [ -f "/data/gsm_sms/config.json" ]; then
    bashio::log.info "GSM SMS configuration created successfully"
    ls -la /data/gsm_sms/config.json
    cat /data/gsm_sms/config.json
else
    bashio::log.error "Failed to create GSM SMS configuration"
fi

# Notify S6 that we're done with this oneshot service
if [ -n "$S6_NOTIFICATION_FD" ]; then
    bashio::log.debug "Notifying S6 that config initialization is complete"
    exec {S6_NOTIFICATION_FD}>&-
fi
